# files/ganglia.generic.cron
*/2 * * * * root /usr/local/bin/network_gmetric.sh
*/2 * * * * root /usr/local/bin/disk_gmetric.sh
